# data/client/borgback.timer
[Unit]
Description=Timer For BorgBack service

[Timer]
OnUnitActiveSec=5m

[Install]
WantedBy=multi-user.target
